# Projects/deployments/daily-fetcher/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor-daily-fetcher.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-daily-fetcher.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-daily-fetcher.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-daily-fetcher.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

;[include]
;files = relative/directory/*.ini


[program:every_day_fetching_Wordpress]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Wordpress -n %%h-every_day_fetching_Wordpress-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Wordpress-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Wordpress-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Wordpress]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q first_fetch.fetching.Wordpress,infrequent.fetching.Wordpress,every_day.fetching.Wordpress -n %%h-infrequent_fetching_Wordpress-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Wordpress-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Wordpress-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Blogspot]
command=./theshelf celery worker --loglevel=INFO -c 5 -Q every_day.fetching.Blogspot --maxtasksperchild 50 -n %%h-every_day_fetching_Blogspot-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Blogspot-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Blogspot-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Blogspot]
command=./theshelf celery worker --loglevel=INFO -c 10 -Q first_fetch.fetching.Blogspot,infrequent.fetching.Blogspot,every_day.fetching.Blogspot --maxtasksperchild 50 -n %%h-infrequent_fetching_Blogspot-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Blogspot-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Blogspot-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Custom]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q every_day.fetching.Custom --maxtasksperchild 50 -n %%h-every_day_fetching_Custom-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Custom-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Custom-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Custom]
command=./theshelf celery worker --loglevel=INFO -c 8 -Q first_fetch.fetching.Custom,infrequent.fetching.Custom,every_day.fetching.Custom --maxtasksperchild 50 -n %%h-infrequent_fetching_Custom-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Custom-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Custom-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Facebook]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Facebook -n %%h-every_day_fetching_Facebook-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Facebook-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Facebook-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Facebook]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q first_fetch.fetching.Facebook,infrequent.fetching.Facebook,every_day.fetching.Facebook -n %%h-infrequent_fetching_Facebook-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Facebook-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Facebook-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Instagram]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Instagram -n %%h-every_day_fetching_Instagram-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Instagram-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Instagram-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Instagram]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q first_fetch.fetching.Instagram,infrequent.fetching.Instagram,every_day.fetching.Instagram -n %%h-infrequent_fetching_Instagram-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Instagram-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Instagram-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Twitter]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Twitter -n %%h-every_day_fetching_Twitter-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Twitter-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Twitter-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Twitter]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q first_fetch.fetching.Twitter,infrequent.fetching.Twitter,every_day.fetching.Twitter -n %%h-infrequent_fetching_Twitter-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Twitter-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Twitter-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Pinterest]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Pinterest --maxtasksperchild 50 -n %%h-every_day_fetching_Pinterest-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Pinterest-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Pinterest-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Pinterest]
command=./theshelf celery worker --loglevel=INFO -c 3 -Q first_fetch.fetching.Pinterest,infrequent.fetching.Pinterest,every_day.fetching.Pinterest --maxtasksperchild 50 -n %%h-infrequent_fetching_Pinterest-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Pinterest-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Pinterest-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Tumblr]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Tumblr -n %%h-every_day_fetching_Tumblr-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Tumblr-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Tumblr-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Tumblr]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q first_fetch.fetching.Tumblr,infrequent.fetching.Tumblr,every_day.fetching.Tumblr -n %%h-infrequent_fetching_Tumblr-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Tumblr-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Tumblr-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Youtube]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Youtube -n %%h-every_day_fetching_Youtube-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Youtube-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Youtube-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:infrequent_fetching_Youtube]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q first_fetch.fetching.Youtube,infrequent.fetching.Youtube,every_day.fetching.Youtube -n %%h-infrequent_fetching_Youtube-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/infrequent_fetching_Youtube-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/infrequent_fetching_Youtube-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_social]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Gplus,every_day.fetching.Bloglovin --maxtasksperchild 50 -n %%h-every_day_fetching_social-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_social-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_social-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:sponsorship_fetching]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q sponsorship_fetching -n %%h-sponsorship_fetching-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/sponsorship_fetching-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/sponsorship_fetching-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:hcheck-rabbitmq-memory]
command=./theshelf python -m servermonitoring.healthchecks check_rabbitmq_memory
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-rabbitmq-memory.out
stderr_logfile=/home/ubuntu/log/hcheck-rabbitmq-memory.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:hcheck-disk]
command=./theshelf python -m servermonitoring.healthchecks check_disk_space
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-disk.out
stderr_logfile=/home/ubuntu/log/hcheck-disk.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:watchdog]
command=./theshelf python -m servermonitoring.watchdog watch
directory=/home/ubuntu/Projects_daily-fetcher
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/watchdog.out
stderr_logfile=/home/ubuntu/log/watchdog.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true
