# Projects/deployments/indepth-fetcher/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor-indepth-fetcher.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-indepth-fetcher.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-indepth-fetcher.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-indepth-fetcher.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

;[include]
;files = relative/directory/*.ini


[program:celeryd-indepth-fetcher-Wordpress]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Wordpress
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Wordpress.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Wordpress.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu

[program:celeryd-indepth-fetcher-Blogspot]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Blogspot
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Blogspot.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Blogspot.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu

[program:celeryd-indepth-fetcher-Facebook]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Facebook
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Facebook.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Facebook.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu

[program:celeryd-indepth-fetcher-Instagram]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Instagram
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Instagram.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Instagram.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu

[program:celeryd-indepth-fetcher-Twitter]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Twitter
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Twitter.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Twitter.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu

[program:celeryd-indepth-fetcher-Pinterest]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Pinterest
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Pinterest.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Pinterest.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu


[program:celeryd-indepth-fetcher-Tumblr]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python manage.py celeryd --loglevel=INFO -c 1 -Q indepth_fetching.Tumblr
directory=/home/ubuntu/Projects_indepth-fetcher/miami_metro
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Tumblr.out
stderr_logfile=/home/ubuntu/log/celeryd-indepth-fetcher-Tumblr.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher",DJANGO_SETTINGS_MODULE="settings"
user=ubuntu

[program:hcheck-disk]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python -m servermonitoring.healthchecks check_disk_space
directory=/home/ubuntu/Projects_indepth-fetcher
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-disk.out
stderr_logfile=/home/ubuntu/log/hcheck-disk.err
autostart=true
autorestart=true
environment=DJANGO_SETTINGS_MODULE="settings",PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher/miami_metro"
user=ubuntu

[program:watchdog]
command=/home/ubuntu/Projects_indepth-fetcher/venv/bin/python -m servermonitoring.watchdog watch
directory=/home/ubuntu/Projects_indepth-fetcher
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/watchdog-indepth-fetcher.out
stderr_logfile=/home/ubuntu/log/watchdog-indepth-fetcher.err
autostart=true
autorestart=true
environment=DJANGO_SETTINGS_MODULE="settings",PYTHONPATH="/home/ubuntu/Projects_indepth-fetcher/miami_metro"
user=ubuntu
