# Projects/deployments/celery-default/files/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file=/tmp/supervisor-celery-default.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

;[inet_http_server]         ; inet (TCP) server disabled by default
;port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[supervisord]
logfile=/tmp/supervisord-celery-default.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-celery-default.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY=value       ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-celery-default.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.


[program:celery]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q celery -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/celery-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/celery-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:celery_google]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q celery -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/celery-google-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/celery-google-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


; WARNING!!! Make sure there is only one beat process and don't run this on many machines.
[program:celery_beat]
command=./theshelf celery beat --loglevel=INFO
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/beat-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/beat-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:email_extraction_high_priority]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q email_extraction_high_priority -n %%h-email_extraction_high_priority-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/email_extraction_high_priority-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/email_extraction_high_priority-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:import_products_from_post_directly_google]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q import_products_from_post_directly -n %%h-import_products_from_post_latest-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/import_products_from_post-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/import_products_from_post-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:new_influencer]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q new_influencer -n %%h-new_influencer-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/new_influencer-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/new_influencer-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:instagram_post_scraper]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q scrape_instagram_posts_new -n %%h-instagram-post-scraper-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/instagram_post_fetcher_%(process_num)d.out
stderr_logfile=/home/ubuntu/log/instagram_post_fetcher_%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

;;; Bookmarklet check not used anymore
;;;[program:hcheck-bookmarklet]
;;;command=./theshelf python -m servermonitoring.healthchecks bookmarklet_processing_test
;;;directory=/home/ubuntu/Projects_celery-default
;;;redirect_stderr = true
;;;stdout_logfile=/home/ubuntu/log/hcheck-bookmarklet.out
;;;stderr_logfile=/home/ubuntu/log/hcheck-bookmarklet.err
;;;autostart=true
;;;autorestart=true
;;;environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
;;;user=ubuntu
;;;killasgroup=true

[program:hcheck-disk]
command=./theshelf python -m servermonitoring.healthchecks check_disk_space
directory=/home/ubuntu/Projects_celery-default
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-disk.out
stderr_logfile=/home/ubuntu/log/hcheck-disk.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:watchdog]
command=./theshelf python -m servermonitoring.watchdog watch
directory=/home/ubuntu/Projects_celery-default
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/watchdog.out
stderr_logfile=/home/ubuntu/log/watchdog.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true


[program:import_products_from_post_directly_google_2]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q import_products_from_post_directly --maxtasksperchild 20 -n %%h-import_products_from_post_directly-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/import_products_from_post_directly-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/import_products_from_post_directly-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu

[program:celery_mandrill_2]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q celery_mandrill_2 -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/celery_mandrill_2-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/celery_mandrill_2-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:bulk_update_campaigns_tracking_stats]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q bulk_update_campaigns_tracking_stats -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/bulk_update_campaigns_tracking_stats-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/bulk_update_campaigns_tracking_stats-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:crawl_contract_influencers]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q crawl_contract_influencers -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/crawl_contract_influencers-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/crawl_contract_influencers-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:crawl_contract_influencers_platforms]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q crawl_contract_influencers_platforms -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/crawl_contract_influencers_platforms-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/crawl_contract_influencers_platforms-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:add_campaign_posts_to_collections]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q add_campaign_posts_to_collections -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/add_campaign_posts_to_collections-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/add_campaign_posts_to_collections-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:perform-social-handle-updates]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q social-handle-updates -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/perform_social_handle_updates-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/perform_social_handle_updates-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:update_bloggers_cache_data]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q update_bloggers_cache_data -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/update_bloggers_cache_data-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/update_bloggers_cache_data-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:fix_pin_date_01]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q fix_pin_date_01 -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/fix_pin_date_01-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/fix_pin_date_01-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:instagram_feed_scraper]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q instagram_feed_scraper -n %%h-instagram-feed-fetcher-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/instagram_feed_fetcher_%(process_num)d.out
stderr_logfile=/home/ubuntu/log/instagram_feed_fetcher_%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:bookmarking]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q bookmarking -n %%h-bookmarking-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_celery-default
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/bookmarking_%(process_num)d.out
stderr_logfile=/home/ubuntu/log/bookmarking_%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true