# Projects/deployments/daily-fetcher-social/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor-daily-fetcher-social.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-daily-fetcher-social.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-daily-fetcher-social.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-daily-fetcher-social.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

;[include]
;files = relative/directory/*.ini

[program:every_day_fetching_Facebook]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Facebook -n %%h-every_day_fetching_Facebook-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Facebook-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Facebook-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:every_day_fetching_Instagram]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Instagram -n %%h-every_day_fetching_Instagram-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Instagram-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Instagram-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:every_day_fetching_Twitter]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Twitter -n %%h-every_day_fetching_Twitter-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Twitter-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Twitter-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:every_day_fetching_Pinterest]
command=./theshelf celery worker --loglevel=INFO -c 5 -Q every_day.fetching.Pinterest --maxtasksperchild 50 -n %%h-every_day_fetching_Pinterest-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Pinterest-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Pinterest-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:every_day_fetching_Tumblr]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Tumblr -n %%h-every_day_fetching_Tumblr-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Tumblr-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Tumblr-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:every_day_fetching_Youtube]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Youtube -n %%h-every_day_fetching_Youtube-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Youtube-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Youtube-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true



[program:every_day_fetching_social]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Gplus,every_day.fetching.Bloglovin --maxtasksperchild 50 -n %%h-every_day_fetching_social-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_social-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_social-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:hcheck-disk]
command=./theshelf python -m servermonitoring.healthchecks check_disk_space
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-disk.out
stderr_logfile=/home/ubuntu/log/hcheck-disk.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:watchdog]
command=./theshelf python -m servermonitoring.watchdog watch
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/watchdog.out
stderr_logfile=/home/ubuntu/log/watchdog.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true


[program:fetch_instagram_post_interactions]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q fetch_instagram_post_interactions --maxtasksperchild 50 -n %%h-fetch_instagram_post_interactions-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/fetch_instagram_post_interactions-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/fetch_instagram_post_interactions-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:customer_uploaded_fetching_Facebook]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Facebook -n %%h-customer_uploaded_fetching_Facebook-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Facebook-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Facebook-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:customer_uploaded_fetching_Instagram]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Instagram -n %%h-customer_uploaded_fetching_Instagram-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Instagram-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Instagram-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:customer_uploaded_fetching_Twitter]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Twitter -n %%h-customer_uploaded_fetching_Twitter-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Twitter-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Twitter-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:customer_uploaded_fetching_Pinterest]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q every_day.fetching.Pinterest --maxtasksperchild 50 -n %%h-customer_uploaded_fetching_Pinterest-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Pinterest-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Pinterest-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:customer_uploaded_fetching_Tumblr]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Tumblr -n %%h-customer_uploaded_fetching_Tumblr-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Tumblr-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Tumblr-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:customer_uploaded_fetching_Youtube]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Youtube -n %%h-customer_uploaded_fetching_Youtube-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Youtube-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Youtube-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:refetch_campaign_posts_Facebook]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Facebook -n %%h-refetch_campaign_posts_Facebook-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Facebook-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Facebook-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:refetch_campaign_posts_Instagram]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Instagram -n %%h-refetch_campaign_posts_Instagram-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Instagram-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Instagram-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:refetch_campaign_posts_Twitter]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Twitter -n %%h-refetch_campaign_posts_Twitter-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Twitter-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Twitter-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:refetch_campaign_posts_Pinterest]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Pinterest --maxtasksperchild 50 -n %%h-refetch_campaign_posts_Pinterest-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Pinterest-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Pinterest-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:refetch_campaign_posts_Tumblr]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Tumblr -n %%h-refetch_campaign_posts_Tumblr-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Tumblr-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Tumblr-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:refetch_campaign_posts_Youtube]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Youtube -n %%h-refetch_campaign_posts_Youtube-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Youtube-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Youtube-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:fix_pin_date_03]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q fix_pin_date_03 -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/fix_pin_date_03-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/fix_pin_date_03-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:discover_new_infs_posts_Instagram]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q discover_new_infs.fetching.Instagram -n %%h-discover_new_infs_Instagram-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-social
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/discover_new_infs_Instagram-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/discover_new_infs_Instagram-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true