# Projects/deployments/rs-daily-fetcher/sysfiles/etc/cron.d/clean_selenium_processes
*/10 * * * *   root    killall -KILL -o 1860s firefox
*/10 * * * *   root    killall -KILL -o 1860s Xvfb
