# Projects/deployments/common/sysfiles/etc/cron.d/clean_selenium_processes
*/10 * * * *   root    killall -KILL -o 1860s firefox
