# Projects/deployments/platform-data-postprocessing/files/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file=/tmp/supervisor-platform-data-postprocessing.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

;[inet_http_server]         ; inet (TCP) server disabled by default
;port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[supervisord]
logfile=/tmp/supervisord-platform-data-postprocessing.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-platform-data-postprocessing.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY=value       ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-platform-data-postprocessing.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:platform_data_postprocessing]
command=./theshelf celery worker --loglevel=INFO -c 12 -Q platform_data_postprocessing --maxtasksperchild 20 -n %%h-platform_data_postprocessing-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/platform_data_postprocessing-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/platform_data_postprocessing-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:platform_data_postprocessing_blocking]
command=./theshelf celery worker --loglevel=INFO -c 6 -Q platform_data_postprocessing_blocking --maxtasksperchild 20 -n %%h-platform-data-postprocessing-blocking%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/platform_data_postprocessing_blocking-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/platform_data_postprocessing_blocking-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:platform_data_content_estimation]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q platform_data_content_estimation -n %%h-platform-data-content-estimation-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/platform_data_content_estimation-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/platform_data_content_estimation-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:platform_extraction]
command=./theshelf celery worker --loglevel=INFO -c 3 -Q platform_extraction --maxtasksperchild 20 -n %%h-platform-extraction-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/platform_extraction-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/platform_extraction-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:platform_extraction_2]
command=./theshelf celery worker --loglevel=INFO -c 3 -Q platform_extraction_2 --maxtasksperchild 20 -n %%h-platform-extraction-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/platform_extraction_2-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/platform_extraction_2-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:blog_visit]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q blog_visit --maxtasksperchild 20 -n %%h-blog-visit-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/blog_visit-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/blog_visit-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:pdimport]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q pdimport --maxtasksperchild 100 -n %%h-pdimport-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/pdimport-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/pdimport-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:post_image]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q post_image_upload_worker --maxtasksperchild 100 -n %%h-postimage-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/postimage-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/postimage-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:hcheck-disk]
command=./theshelf python -m servermonitoring.healthchecks check_disk_space
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-disk-platform-data-postprocessing.out
stderr_logfile=/home/ubuntu/log/hcheck-disk-platform-data-postprocessing.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:watchdog]
command=./theshelf python -m servermonitoring.watchdog watch
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/watchdog.out
stderr_logfile=/home/ubuntu/log/watchdog.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google_pgbouncer",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true


[program:post_categorization]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q post_categorization -n post_cat%(process_num)d.%%h --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/post_categorization-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/post_categorization-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:newinfluencer]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q new_influencer -n %%h-fetcher-new-influencer-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/newinfluencer_%(process_num)d.out
stderr_logfile=/home/ubuntu/log/newinfluencer_%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:instagram_post_scraper]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q scrape_instagram_posts_new -n %%h-instagram-post-scraper-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/instagram_post_fetcher_%(process_num)d.out
stderr_logfile=/home/ubuntu/log/instagram_post_fetcher_%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true


[program:import_products_from_post_directly]
command=./theshelf celery worker --loglevel=INFO -c 15 -Q import_products_from_post_directly --maxtasksperchild 20 -n %%h-import_products_from_post_directly-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/import_products_from_post_directly-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/import_products_from_post_directly-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu

[program:instagram_feed_scraper]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q instagram_feed_scraper -n %%h-instagram-feed-fetcher-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_platform-data-postprocessing
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/instagram_feed_fetcher_%(process_num)d.out
stderr_logfile=/home/ubuntu/log/instagram_feed_fetcher_%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true