# Projects/deployments/daily-fetcher-blogs/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor-daily-fetcher-blogs.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-daily-fetcher-blogs.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-daily-fetcher-blogs.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-daily-fetcher-blogs.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

;[include]
;files = relative/directory/*.ini


[program:every_day_fetching_Wordpress]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q every_day.fetching.Wordpress -n %%h-every_day_fetching_Wordpress-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Wordpress-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Wordpress-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Blogspot]
command=./theshelf celery worker --loglevel=INFO -c 5 -Q every_day.fetching.Blogspot --maxtasksperchild 50 -n %%h-every_day_fetching_Blogspot-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Blogspot-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Blogspot-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Custom]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q every_day.fetching.Custom --maxtasksperchild 50 -n %%h-every_day_fetching_Custom-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Custom-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Custom-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:every_day_fetching_Squarespace]
command=./theshelf celery worker --loglevel=INFO -c 4 -Q every_day.fetching.Squarespace --maxtasksperchild 50 -n %%h-every_day_fetching_Squarespace-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/every_day_fetching_Squarespace-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/every_day_fetching_Squarespace-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:sponsorship_fetching]
command=./theshelf celery worker --loglevel=INFO -c 2 -Q sponsorship_fetching -n %%h-sponsorship_fetching-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/sponsorship_fetching-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/sponsorship_fetching-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:hcheck-rabbitmq-memory]
command=./theshelf python -m servermonitoring.healthchecks check_rabbitmq_memory
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-rabbitmq-memory.out
stderr_logfile=/home/ubuntu/log/hcheck-rabbitmq-memory.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:hcheck-disk]
command=./theshelf python -m servermonitoring.healthchecks check_disk_space
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/hcheck-disk.out
stderr_logfile=/home/ubuntu/log/hcheck-disk.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:watchdog]
command=./theshelf python -m servermonitoring.watchdog watch
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/watchdog.out
stderr_logfile=/home/ubuntu/log/watchdog.err
autostart=true
autorestart=true
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
stopasgroup=true
killasgroup=true

[program:import_products_from_post_directly]
command=./theshelf celery worker --loglevel=INFO -c 5 -Q import_products_from_post_directly --maxtasksperchild 20 -n %%h-import_products_from_post_directly-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr = true
stdout_logfile=/home/ubuntu/log/import_products_from_post_directly-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/import_products_from_post_directly-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu

[program:customer_uploaded_fetching_Wordpress]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Wordpress -n %%h-customer_uploaded_fetching_Wordpress-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Wordpress-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Wordpress-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:customer_uploaded_fetching_Blogspot]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Blogspot --maxtasksperchild 50 -n %%h-customer_uploaded_fetching_Blogspot-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Blogspot-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Blogspot-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:customer_uploaded_fetching_Custom]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Custom --maxtasksperchild 50 -n %%h-customer_uploaded_fetching_Custom-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Custom-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Custom-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:customer_uploaded_fetching_Squarespace]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q customer_uploaded.fetching.Squarespace --maxtasksperchild 50 -n %%h-customer_uploaded_fetching_Squarespace-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/customer_uploaded_fetching_Squarespace-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/customer_uploaded_fetching_Squarespace-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:refetch_campaign_posts_Wordpress]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Wordpress -n %%h-refetch_campaign_posts_Wordpress-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Wordpress-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Wordpress-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:refetch_campaign_posts_Blogspot]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Blogspot --maxtasksperchild 50 -n %%h-refetch_campaign_posts_Blogspot-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Blogspot-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Blogspot-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:refetch_campaign_posts_Custom]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Custom --maxtasksperchild 50 -n %%h-refetch_campaign_posts_Custom-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Custom-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Custom-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:refetch_campaign_posts_Squarespace]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q refetch_campaign_posts.Squarespace --maxtasksperchild 50 -n %%h-refetch_campaign_posts_Squarespace-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/refetch_campaign_posts_Squarespace-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/refetch_campaign_posts_Squarespace-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true

[program:fix_pin_date_02]
command=./theshelf celery worker --loglevel=INFO -c 1 -Q fix_pin_date_02 -n %%h-celery-%(process_num)d --without-mingle -O fair
numprocs=1
process_name=%(program_name)s_%(process_num)d
directory=/home/ubuntu/Projects_daily-fetcher-blogs
redirect_stderr=true
stdout_logfile=/home/ubuntu/log/fix_pin_date_02-%(process_num)d.out
stderr_logfile=/home/ubuntu/log/fix_pin_date_02-%(process_num)d.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
environment=HOME="/home/ubuntu",E="google",XVFB="True",DISPLAY=":1"
user=ubuntu
killasgroup=true
